# repo/crontab
@reboot /home/pi/kplex-monitor/run/boot_test.sh

